# templates/dwh-etl.cron
SHELL=/bin/sh
PATH=/usr/local/bin:/bin:/usr/bin

#0 2 * * *	${app.user}	. $HOME/.profile && $DWH_HOME/bin/dwh-update.sh
#0 2 * * *	${app.user}	. $HOME/.profile && $DWH_HOME/bin/psa-update.sh
#0 2 * * *	${app.user}	. $HOME/.profile && $DWH_HOME/bin/dsa-update.sh

0 2 * * *	su - ${app.user} -c "${app.rootDir}/${app.user}/bin/dwh-update.sh" -s /bin/bash
